# other_config/mount/mnt-workstation.automount
[Unit]
Description=Automount synology workstation

[Automount]
Where=/mnt/workstation
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
