# other_config/mount/run-media-afa-Synology_WorkStation.automount
[Unit]
Description=Automount synology workstation

[Automount]
Where=/run/media/afa/Synology_WorkStation/
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
