# other_config/mount/run-media-afa-Synology_FTP.automount
[Unit]
Description=Automount synology FTP Server

[Automount]
Where=/run/media/afa/Synology_FTP/
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
